# deploy/com.fuling.ops-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS LaunchAgent for the standing health monitor (ops_monitor), READ-ONLY reconcilers.
     Execs the python binary DIRECTLY (it lives outside ~/Downloads, so launchd can exec it) — so the
     ONLY thing needing Full Disk Access is this python, not /bin/bash. (The bash wrapper
     deploy/run_ops_monitor.sh is for the Linux/crontab path instead.)

     Prereq on this Mac (repo is under ~/Downloads, which macOS TCC blocks for launchd):
       System Settings → Privacy & Security → Full Disk Access → "+" → ⌘⇧G →
       /usr/bin/python3   → enable the toggle.
     (We use /usr/bin/python3 — arm64-native, has the prod deps in user-site — instead of the conda
      python3.10, which is an x86_64 binary the FDA picker greys out on Apple Silicon.)
     Install:  cp deploy/com.fuling.ops-monitor.plist ~/Library/LaunchAgents/
               launchctl load -w ~/Library/LaunchAgents/com.fuling.ops-monitor.plist
     Test now: launchctl kickstart gui/$(id -u)/com.fuling.ops-monitor   (then: launchctl list | grep fuling)
     Laptop sleeps → fires at 02:30 if awake, else once on next wake. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.fuling.ops-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>opensearch_pipeline.ops_monitor</string>
    <string>--only</string>
    <string>reconcile_ha3</string>
    <string>reconcile_oss</string>
    <!-- reconcile_raw intentionally NOT in the daily set: it exit-2's on the 1 known benign raw_key
         gap (DOC_ADMIN_…14D1C1), which would train RED-blindness. Daily = recall-critical ha3+oss
         (green-when-healthy → RED means a real new problem). Run reconcile_raw ad-hoc:
           RAG_ENV=prod_ro /usr/bin/python3 -m opensearch_pipeline.reconcile --check raw
         (or add a separate weekly agent once that 1 doc's raw is triaged). -->
  </array>
  <key>WorkingDirectory</key><string>/Users/laijunchen/Downloads/opensearch-rag-pipeline</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RAG_ENV</key><string>prod_ro</string>
    <key>PYTHONPATH</key><string>/Users/laijunchen/Downloads/opensearch-rag-pipeline</string>
    <!-- set these to deliver OBS-4 alerts on drift (else the alert path is a logged no-op):
    <key>RAG_OPS_ALERT_WEBHOOK</key><string>https://oapi.dingtalk.com/robot/send?access_token=…</string>
    <key>RAG_OPS_ALERT_SECRET</key><string>SEC…</string> -->
  </dict>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>30</integer></dict>
  <key>RunAtLoad</key><false/>
  <!-- logs MUST live outside ~/Downloads: launchd opens StandardOut/ErrorPath as launchd (not the
       FDA'd python), and macOS TCC blocks launchd from ~/Downloads → EX_CONFIG (78) spawn failure.
       ~/Library/Logs is the conventional, non-TCC-restricted location. -->
  <key>StandardOutPath</key><string>/Users/laijunchen/Library/Logs/fuling-ops-monitor.out</string>
  <key>StandardErrorPath</key><string>/Users/laijunchen/Library/Logs/fuling-ops-monitor.err</string>
</dict>
</plist>
